# .org/config/manjaro_startup/autostart/Flameshot.desktop
[Desktop Entry]
Name=flameshot
Icon=flameshot
Exec=$EXEC_PATH
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
